# dist/linux/runit/armory_utxsvr/run
#!/bin/sh

USER=ubuntu
export USER_HOME=/home/ubuntu

#sv start bitcoind || exit 1
exec 2>&1
exec su -s /bin/bash -c 'cd ${USER_HOME} && /usr/local/bin/armory_utxsvr' ${USER}

