# dist/linux/runit/bitcoind/run
#!/bin/sh

USER=ubuntu
USER_HOME=/home/ubuntu

exec chpst -u ${USER} myriadcoind -daemon=0 -datadir=${USER_HOME}/.myriadcoin
